# supervisord.conf
[supervisord]
nodaemon=true

[program:rbot]
command=bash /app/sh_client_bot.sh 
stdout_logfile=/var/log/log_r_client.log.log
stderr_logfile=/var/log/log_r_client.log_err.log
